# conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:confluence]
command=/opt/atlassian/jira/bin/start-jira.sh -fg jirauser
autorestart=false


